# backend/leonbos.service
[Unit]
Description=LeonBos Backend Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/leonbos
Environment=NODE_ENV=production
Environment=LEONBOS_DATA_DIR=/var/lib/leonbos
ExecStart=/bin/bash /opt/leonbos/backend/start.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
